# deploy/systemd/ble-command-gateway.service
[Unit]
Description=BLE Command Gateway Server
After=bluetooth.service NetworkManager.service
Wants=bluetooth.service NetworkManager.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/ble-command-gateway
ExecStart=/opt/ble-command-gateway/target/release/server
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
